# systemd/kiosk-update.timer
[Unit]
Description=KNIGER Kiosk Remote-Update Timer (nächtliches Wartungsfenster)
Documentation=https://github.com/behrndt/kniger-kiosk

[Timer]
# Bei aktivem Overlay-FS löst ein Update einen Reboot-Zyklus aus (Overlay aus →
# pull → Overlay an → reboot). Deshalb NICHT stündlich, sondern nachts im
# Wartungsfenster, wenn der Kiosk-Bildschirm ohnehin aus ist.
OnCalendar=*-*-* 03:30:00
# Jitter, damit mehrere Kiosks nicht exakt gleichzeitig pullen/rebooten
RandomizedDelaySec=20min
Persistent=false
Unit=kiosk-update.service
AccuracySec=5min

[Install]
WantedBy=timers.target
